# scripts/jarvis/nex-weekly-bench.timer
[Unit]
Description=Run full nex-code benchmark weekly (Saturday 04:00)

[Timer]
OnCalendar=Sat *-*-* 04:00:00
Persistent=true

[Install]
WantedBy=timers.target
